# ip-proxy/ipproxy.ini
[program:ipproxy]
directory = /data/app
command = /usr/java/latest/bin/java -Dserver.port=8808
    -Dlogging.config=/data/app/logback-spring-prod.xml
    -Dspring.profiles.active=prod
    -jar /data/app/ip-proxy-0.0.1.jar
stdout_logfile = /data/app/ipproxy-stdout.log
stderr_logfile = /data/app/ipproxy-stderr.log
startsecs=30
startretries=1

